# kafka/kafka.conf
[program:kafka]
command=/usr/bin/kafka.sh
autostart=true
autorestart=true